# luci-app-simple2fa/root/etc/init.d/simple2fa
#!/bin/sh /etc/rc.common
# Simple2FA Debug Build - init.d script
# 所有操作都会记录到系统日志，使用 logread | grep simple2fa 查看

START=99
STOP=10

# 目标文件
CGI_TARGET="/www/cgi-bin/luci"
SYSAUTH_TARGET="/usr/lib/lua/luci/view/sysauth.htm"

# 源文件 (我们的包)
CGI_SOURCE="/usr/share/luci-app-simple2fa/luci"
SYSAUTH_SOURCE="/usr/share/luci-app-simple2fa/sysauth.htm"

log() {
	logger -t simple2fa "[init.d] $1"
	echo "[simple2fa] $1"
}

start() {
	log "========== START 被调用 =========="
	
	local enabled=$(uci -q get simple2fa.global.enabled)
	log "配置读取: simple2fa.global.enabled = '$enabled'"
	
	local secret=$(uci -q get simple2fa.global.secret)
	if [ -n "$secret" ]; then
		log "密钥状态: 已配置 (长度=${#secret})"
	else
		log "密钥状态: 未配置!"
	fi

	if [ "$enabled" = "1" ]; then
		log "开关状态: 启用 -> 开始替换文件"
		
		# === 1. 备份 CGI ===
		if [ -e "$CGI_TARGET" ]; then
			if [ ! -e "${CGI_TARGET}.bak" ]; then
				if ! grep -q "Simple2FA-Debug" "$CGI_TARGET" 2>/dev/null; then
					cp -af "$CGI_TARGET" "${CGI_TARGET}.bak"
					log "CGI 备份: 成功创建 ${CGI_TARGET}.bak"
				else
					log "CGI 备份: 跳过 (已是我们的脚本)"
				fi
			else
				log "CGI 备份: 跳过 (备份已存在)"
			fi
		else
			log "CGI 备份: 错误! 目标文件不存在: $CGI_TARGET"
		fi
		
		# === 2. 备份 sysauth.htm ===
		if [ -e "$SYSAUTH_TARGET" ]; then
			if [ ! -e "${SYSAUTH_TARGET}.bak" ]; then
				cp -af "$SYSAUTH_TARGET" "${SYSAUTH_TARGET}.bak"
				log "模板备份: 成功创建 ${SYSAUTH_TARGET}.bak"
			else
				log "模板备份: 跳过 (备份已存在)"
			fi
		else
			log "模板备份: 警告! 目标文件不存在: $SYSAUTH_TARGET"
		fi
		
		# === 3. 安装我们的 CGI ===
		if [ -f "$CGI_SOURCE" ]; then
			rm -f "$CGI_TARGET"
			cp -f "$CGI_SOURCE" "$CGI_TARGET"
			chmod +x "$CGI_TARGET"
			log "CGI 安装: 成功 ($CGI_SOURCE -> $CGI_TARGET)"
		else
			log "CGI 安装: 错误! 源文件不存在: $CGI_SOURCE"
		fi
		
		# === 4. 安装我们的 sysauth.htm ===
		if [ -f "$SYSAUTH_SOURCE" ]; then
			rm -f "$SYSAUTH_TARGET"
			cp -f "$SYSAUTH_SOURCE" "$SYSAUTH_TARGET"
			log "模板安装: 成功 ($SYSAUTH_SOURCE -> $SYSAUTH_TARGET)"
		else
			log "模板安装: 错误! 源文件不存在: $SYSAUTH_SOURCE"
		fi
		
		log "========== START 完成 (2FA 已激活) =========="
	else
		log "开关状态: 禁用 -> 调用 stop 恢复原始文件"
		stop
		log "========== START 完成 (2FA 未激活) =========="
	fi
}

stop() {
	log "========== STOP 被调用 =========="
	
	# === 1. 恢复 CGI ===
	if [ -f "${CGI_TARGET}.bak" ]; then
		rm -f "$CGI_TARGET"
		cp -af "${CGI_TARGET}.bak" "$CGI_TARGET"
		rm -f "${CGI_TARGET}.bak"
		chmod +x "$CGI_TARGET"
		log "CGI 恢复: 成功 (从 .bak 恢复)"
	else
		log "CGI 恢复: 跳过 (无备份文件)"
	fi

	# === 2. 恢复 sysauth.htm ===
	if [ -f "${SYSAUTH_TARGET}.bak" ]; then
		rm -f "$SYSAUTH_TARGET"
		cp -af "${SYSAUTH_TARGET}.bak" "$SYSAUTH_TARGET"
		rm -f "${SYSAUTH_TARGET}.bak"
		log "模板恢复: 成功 (从 .bak 恢复)"
	else
		log "模板恢复: 跳过 (无备份文件)"
	fi
	
	log "========== STOP 完成 =========="
}

restart() {
	log "========== RESTART 被调用 =========="
	stop
	start
}
